# config/rc.local
#!/bin/sh -e

/home/pi/TeleHomeGuard/run.sh || exit 1

exit 0